# systemd/backup-k8s-resources.timer
[Unit]
Description=Timer for Kubernetes Resources Backup

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
