# systemd/backup-k8s-resources.service
[Unit]
Description=Kubernetes Resources Backup

[Service]
User=root
Group=root
Type=oneshot
ExecStart=/usr/local/bin/backup-k8s-resources.sh

[Install]
WantedBy=multi-user.target
